# agent_mitmproxy/docker/images/agent_mitmproxy/supervisord.conf
; ═══════════════════════════════════════════════════════════════════════════════
; Agent Mitmproxy — supervisord config
;
; Two sibling programs. supervisord is PID 1; it forwards signals and
; restarts children if they crash. Both programs log to the container
; stdout so `docker logs` / fluent-bit sees everything.
; ═══════════════════════════════════════════════════════════════════════════════

[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:mitmweb]
command=mitmweb
        --proxyauth "%(ENV_AGENT_MITMPROXY__PROXY_AUTH_USER)s:%(ENV_AGENT_MITMPROXY__PROXY_AUTH_PASS)s"
        --web-host 127.0.0.1
        --web-port 8081
        --listen-port 8080
        --set block_global=false
        --set confdir=/root/.mitmproxy
        -s /app/agent_mitmproxy/addons/addon_registry.py
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:admin_api]
command=uvicorn agent_mitmproxy.fast_api.app:app --host 0.0.0.0 --port 8000
directory=/app
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
