# sg_compute_specs/mitmproxy/docker/images/agent_mitmproxy/supervisord.conf
; ═══════════════════════════════════════════════════════════════════════════════
; Agent Mitmproxy — supervisord config
;
; Two sibling programs. supervisord is PID 1; it forwards signals and
; restarts children if they crash. Both programs log to the container
; stdout so `docker logs` / fluent-bit sees everything.
; ═══════════════════════════════════════════════════════════════════════════════

[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:mitmweb]
command=/bin/sh /tmp/run_mitmweb.sh
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:admin_api]
command=uvicorn sg_compute_specs.mitmproxy.api.app:app --host 0.0.0.0 --port 8000
directory=/app
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
